# larch/profiles/kde/rootoverlay/etc/profile.d/kde_no_ipv6.sh
# Konqueror seems to access the web very slowly without this
export KDE_NO_IPV6=true
