# larch-live/etc/profile.d/larch-live.sh
export PATH=${PATH}:/opt/larch-live/run
